# src/software/jetson_nano/systemd/redis_watcher.path
[Unit]
Description=Monitor for changes in the redis script
After=thunderbots.service

[Path]
#todo add binary path
PathChanged=${%i_path}
Unit=watcher@%i.service
