# src/software/jetson_nano/systemd/redis_watcher.path
[Unit]
Description=Monitor for changes in the redis script
After=thunderbots.service

[Path]
#todo add binary path (#2443)
PathChanged=/home/robot/example.txt
Unit=watcher@redis.service

[Install]
# This unit should be run at startup
WantedBy=multi-user.target
