# src/software/jetson_nano/systemd/tftp_watcher.path
[Unit]
Description=Monitor for changes in the tftp script
After=thunderbots.service

[Path]
PathChanged=/home/robot/thunderbots_binaries/tftp/receiver/receiver.runfiles/__main__/software/jetson_nano/tftp/receiver.py
Unit=watcher@tftp.service

[Install]
# This unit should be run at startup
WantedBy=multi-user.target
