# src/software/jetson_nano/systemd/thunderloop_watcher.path
[Unit]
Description=Monitor the file for changes in the thunderloop binary
After=thunderbots.service

[Path]
#todo add binary path
PathChanged=${%i_path}
Unit=watcher@%i.service
