# src/software/jetson_nano/systemd/thunderloop_watcher.path
[Unit]
Description=Monitor the file for changes in the thunderloop binary
After=thunderbots.service

[Path]
PathChanged=/home/robot/thunderbots_binaries/thunderloop_main
Unit=watcher@hunderloop.service

[Install]
# This unit should be run at startup
WantedBy=multi-user.target
