# src/software/jetson_nano/systemd/ui_watcher.path
[Unit]
Description=Monitor for changes in the announcement script
After=thunderbots.service

[Path]
#todo add binary path
PathChanged=${%i_path}
Unit=watcher@%announcement.service

[Install]
# This unit should be run at startup
WantedBy=multi-user.target